# files/profile.d/03-phpcs-aliases.sh
# Aliases for phpcs and phpcbf
alias phpcs="php /usr/local/bin/phpcs/bin/phpcs --standard=PSR2 \$1"
alias phpcbf="php /usr/local/bin/phpcs/bin/phpcbf --standard=PSR2 \$1"
alias phpcs-wp="php /usr/local/bin/phpcs/bin/phpcs --standard=WordPress \$1"
alias phpcbf-wp="php /usr/local/bin/phpcs/bin/phpcbf --standard=WordPress \$1"
